# 2rtk.service
[Unit]
Description=2RTK NTRIP Caster Service
After=network.target

[Service]
Type=simple
User=2rtk
Group=2rtk
WorkingDirectory=/opt/2rtk
ExecStart=/usr/bin/python3 /opt/2rtk/2rtk.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=2rtk

# 安全设置
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ReadWritePaths=/var/log/2rtk

[Install]
WantedBy=multi-user.target